# target/linux/bcm27xx/base-files/etc/uci-defaults/90_rpi-4b-set-wifi
. /lib/functions.sh

case "$(board_name)" in
	raspberrypi,4-model-b)
		uci -q set wireless.radio0.htmode=VHT40
		uci -q set wireless.radio0.channel=36
		uci -q commit wireless
		;;
esac

exit 0
